# mount/install_ssh.service
[Unit]
Description=Install SSH server.

[Service]
Type=oneshot
ExecStart=/home/pi/upload/install_ssh.sh

[Install]
WantedBy=multi-user.target