# mount/rc.local
#!/bin/sh -e

/home/pi/upload/install_ssh.sh

exit 0